# proxmox-config/templates/ansible/templates/growroot.service.j2
[Unit]
Description=Grow root partition

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart="{{ growroot_bin_location }}"

[Install]
WantedBy=sysinit.target
